# home/plasma/.config/autostart/swhkd.desktop
[Desktop Entry]
Exec=swhkd
Icon=
Name=swhkd
Path=
Terminal=False
Type=Application
