# home/plasma/.config/autostart/swhkd.desktop
[Desktop Entry]
Exec=~/Scripts/run_swhkd.sh
Icon=
Name=swhkd
Path=
Terminal=False
Type=Application
